# templates/noblank.service.j2
[Unit]
After=x11.service

[Service]
ExecStart=/usr/local/bin/noblank.sh

[Install]
WantedBy=multi-user.target
