# templates/x11.service.j2
[Unit]
After=systemd-user-sessions.service

[Service]
ExecStart=X

[Install]
WantedBy=multi-user.target
